# docker/rootfs/etc/services.d/hue-fireplace/run
#!/usr/bin/with-contenv ash
# shellcheck shell=ash

set -eo pipefail

cd /hue-fireplace
exec node ./main.js
